# clash/clash.service
[Unit]
Description=clash systemd service.

[Service]
Type=simple
ExecStart=/usr/bin/clash -d /etc/clash

[Install]
WantedBy=multi-user.target
